# misc/otel-upload.service
[Unit]
Description=BAIIAB OpenTelemetry Upload Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/baiiab
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/pi/baiiab/helpers/otel_upload_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
